# Package/Library/LaunchDaemons/io.macadmins.Outset.boot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Program</key>
	<array>
		<string>/usr/local/outset/Outset.app/Contents/MacOS/Outset</string>
	</array>
	<key>AssociatedBundleIdentifiers</key>
	<string>io.macadmins.Outset</string>
	<key>Label</key>
	<string>io.macadmins.Outset.boot</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/outset/Outset.app/Contents/MacOS/Outset</string>
		<string>--boot</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
