# Package/Library/LaunchDaemons/io.macadmins.Outset.login-privileged.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Program</key>
	<string>/usr/local/outset/Outset.app/Contents/MacOS/Outset</string>
	<key>AssociatedBundleIdentifiers</key>
	<string>io.macadmins.Outset</string>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/private/tmp/.io.macadmins.outset.login-privileged.launchd</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>io.macadmins.Outset.login-privileged</string>
	<key>OnDemand</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/outset/Outset.app/Contents/MacOS/Outset</string>
		<string>--login-privileged</string>
	</array>
</dict>
</plist>
